# drive-monitor.service
[Unit]
Description=Drive Monitor Service
After=network.target

[Service]
Type=simple
WorkingDirectory=INSTALL_DIR
ExecStart=/bin/bash INSTALL_DIR/drive_monitor.sh
Environment="PATH=INSTALL_DIR/.venv/bin:$PATH"
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target