# drive-monitor.service
[Unit]
Description=Drive Monitor Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/syncer
ExecStart=/bin/bash %h/syncer/drive_monitor.sh
Environment="PATH=%h/syncer/.venv/bin:$PATH"
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target